# deploy/services/nectarcollector.service
[Unit]
Description=NectarCollector CDR Capture Service
Documentation=https://github.com/911logix/nectarcollector
After=network.target nats-server.service
Wants=nats-server.service

[Service]
Type=simple
User=nectarcollector
Group=dialout
ExecStart=/usr/local/bin/nectarcollector -config /etc/nectarcollector/config.json
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/nectarcollector
PrivateTmp=true

[Install]
WantedBy=multi-user.target
